# 90-garmin.rules
ACTION=="add", \
  KERNEL=="sd*", \
  ATTRS{idVendor}=="*", \
  ATTRS{idProduct}=="*", \
  RUN+="/home/augustus/code/fit-upload/upload.sh"
