# 90-garmin.rules
KERNEL=="sd*", ATTRS{idVendor}=="091e", ATTRS{idProduct}=="086d", RUN+="/home/augustus/code/upload.sh"
