# 90-garmin.rules
ACTION=="add", \
  KERNEL=="sd*", \
  ATTRS{idVendor}=="091e", \
  ATTRS{idProduct}=="086d", \
  RUN+="/home/augustus/code/fit-upload/upload.sh"
